# .config/systemd/user/mbsync.service
[Unit]
Description=Mail synchronization service

[Service]
Type=oneshot
ExecStart=mbsync -Va
ExecStartPost=notmuch new
