# .config/systemd/user/mbsync.service
[Unit]
Description=Mail synchronization service

[Service]
Type=oneshot
ExecStart=%h/.bin/sync-mailbox.sh
