# .config/autostart/01-xsettingsd.desktop
[Desktop Entry]
Name=xsettingd
GenericName=A lightweight xsettings daemon
Exec=xsettingsd
Type=Application
Terminal=false
StartupNotify=false
