# .config/autostart/clipit-startup.desktop
[Desktop Entry]
Name=ClipIt
Comment=Clipboard Manager
Icon=clipit-trayicon
Exec=clipit -d
Terminal=false
Type=Application
