# .config/autostart/compton.desktop
[Desktop Entry]
Name=Compton
GenericName=A compositor for X
Exec=compton -b
Terminal=false
Type=Application
StartupNotify=false
