# .config/autostart/gromit-mpx.desktop
[Desktop Entry]
Name=Gromit-MPX
GenericName=On Screen Drawing Tool
Exec=gromit-mpx
Type=Application
StartupNotify=false

