# .config/autostart/nitrogen.desktop
[Desktop Entry]
Name=Nitrogen
Comment=Desktop background setter for X Window
Type=Application
Exec=nitrogen --restore
StartupNotify=false
